# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/massgo/repos/flask-uwsgi-nginx/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
